# tequila/tequila.crontab
0 0 * * * /opt/tequila/bin/cronjob_vacation